# deploy/yt2ptd.service
[Unit]
Description=yt2pt Daemon — YouTube to PeerTube sync service
Documentation=https://github.com/jloup83/yt2pt
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=yt2pt
Group=yt2pt
WorkingDirectory=/var/lib/yt2pt
ExecStart=/usr/local/bin/yt2ptd
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=yt2ptd

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/yt2pt /var/log/yt2pt /etc/yt2pt
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
RestrictNamespaces=true
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictRealtime=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6

[Install]
WantedBy=multi-user.target
